# templates/systemd/docker-registry-purger.service.j2
[Unit]
Description=Performs container registry purging

[Service]
Type=oneshot
ExecStart={{ docker_registry_purger_bin_path }}/purge
